# supervisor/supervisord.conf
[supervisord]
logfile=/dev/null
pidfile=/home/user/supervisord.pid
nodaemon=true
strip_ansi=false
childlogdir=/home/user/logs

[program:syncd]
command=/bin/bash -lc 'if [ -n "${GITHUB_REPO:-}" ] && [ -n "${GITHUB_PAT:-}" ]; then exec /home/user/.venv/bin/python -m sync; else echo "[syncd] GITHUB_REPO/GITHUB_PAT not set; sync disabled"; exec sleep infinity; fi'
directory=/home/user
autostart=true
autorestart=true
startsecs=1
priority=10
stopsignal=TERM
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:nginx]
command=/bin/bash -lc '/usr/local/openresty/bin/openresty -p /home/user/nginx -c /home/user/nginx/nginx.conf -g "daemon off;"'
directory=/home/user/nginx
autostart=true
autorestart=true
priority=5
stopsignal=TERM
stdout_logfile=/home/user/logs/nginx.supervisor.out
stderr_logfile=/home/user/logs/nginx.supervisor.err
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:codex2api]
command=/bin/bash -lc '/home/user/scripts/run-codex2api.sh'
directory=/home/user
user=1000
autostart=true
autorestart=true
startsecs=2
priority=20
stopsignal=TERM
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:sqlite-backup]
command=/bin/bash -lc '/home/user/scripts/backup-sqlite.sh'
directory=/home/user
user=1000
autostart=true
autorestart=true
startsecs=1
priority=25
stopsignal=TERM
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:filebrowser]
command=/bin/bash -lc '/home/user/scripts/wait-for-sync.sh && /home/user/filebrowser --address 0.0.0.0 --port 8888 --root / --baseurl /filebrowser --database /home/user/filebrowser-data/filebrowser.db'
directory=/home/user
autostart=true
autorestart=true
priority=40
stopsignal=TERM
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:gotty]
command=/bin/bash -lc '/home/user/gotty --address 127.0.0.1 --port 8081 --path /t --credential ${GOTTY_USERNAME:-admin}:${GOTTY_PASSWORD:-adminadminadmin} -w bash'
directory=/home/user
autostart=true
autorestart=true
priority=45
stopsignal=TERM
stdout_logfile=/home/user/logs/gotty.supervisor.out
stderr_logfile=/home/user/logs/gotty.supervisor.err
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
